# src/service/guise.service
# should be placed in:
# /usr/local/lib/systemd/system
# check paths with:
# systemd-analyze unit-paths

[Unit]
Description=Guise User Login
After=network.target
Wants=network.target

[Service]
ExecStart=/usr/local/sbin/guised
ExecReload=/usr/local/sbin/guised
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
